# backend/supervisor/laravel-queue-core.conf
[program:laravel-queue-default]
command=/usr/local/bin/php artisan queue:work --queue=default --sleep=1 --tries=3 --timeout=90 --max-time=3600 --memory=128 --backoff=3,10,30
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=5
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:laravel-queue-log-rotation]
command=/usr/local/bin/php artisan queue:work --queue=log-rotation --sleep=15 --tries=3 --timeout=120 --max-time=3600 --memory=128
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=90
stopsignal=TERM
priority=30
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:laravel-queue-payments]
command=/usr/local/bin/php artisan queue:work --queue=payments --sleep=1 --tries=3 --timeout=120 --max-time=1800 --memory=128 --backoff=5,15,60
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=120
stopsignal=TERM
priority=5
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:laravel-queue-payment-checks]
command=/usr/local/bin/php artisan queue:work --queue=payment-checks --sleep=1 --tries=3 --timeout=90 --max-time=1800 --memory=128 --backoff=5,15,60
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=5
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:laravel-queue-hotspot-sms]
command=/usr/local/bin/php artisan queue:work --queue=hotspot-sms --sleep=1 --tries=3 --timeout=45 --max-time=1800 --memory=128 --backoff=3,10,30
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=5
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:laravel-queue-hotspot-sessions]
command=/usr/local/bin/php artisan queue:work --queue=hotspot-sessions --sleep=1 --tries=3 --timeout=60 --max-time=1800 --memory=128 --backoff=3,10,30
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=10
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:laravel-queue-hotspot-accounting]
command=/usr/local/bin/php artisan queue:work --queue=hotspot-accounting --sleep=1 --tries=2 --timeout=120 --max-time=3600 --memory=128 --backoff=10,30,60
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=15
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:laravel-queue-service-control]
command=/usr/local/bin/php artisan queue:work --queue=service-control --sleep=1 --tries=3 --timeout=90 --max-time=1800 --memory=128 --backoff=3,10,30
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=5
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:laravel-queue-security]
command=/usr/local/bin/php artisan queue:work --queue=security --sleep=1 --tries=3 --timeout=60 --max-time=3600 --memory=128 --backoff=3,10,30
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=5
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:laravel-queue-monitoring]
command=/usr/local/bin/php artisan queue:work --queue=monitoring --sleep=1 --tries=3 --timeout=300 --max-time=3600 --memory=128 --backoff=10,30,60
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=5
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:laravel-queue-tenant-management]
command=/usr/local/bin/php artisan queue:work --queue=tenant-management --sleep=1 --tries=3 --timeout=180 --max-time=3600 --memory=256 --backoff=5,15,30
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=90
stopsignal=TERM
priority=5
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:laravel-queue-user-management]
command=/usr/local/bin/php artisan queue:work --queue=user-management --sleep=1 --tries=3 --timeout=120 --max-time=3600 --memory=256 --backoff=5,15,30
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=5
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:laravel-queue-auth-tracking]
command=/usr/local/bin/php artisan queue:work --queue=auth-tracking --sleep=1 --tries=3 --timeout=60 --max-time=3600 --memory=128 --backoff=3,10,30
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=10
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:laravel-queue-packages]
command=/usr/local/bin/php artisan queue:work --queue=packages --sleep=1 --tries=3 --timeout=90 --max-time=3600 --memory=128 --backoff=3,10,30
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=5
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[group:laravel-queues-core]
programs=laravel-queue-default,laravel-queue-log-rotation,laravel-queue-payments,laravel-queue-payment-checks,laravel-queue-hotspot-sms,laravel-queue-hotspot-sessions,laravel-queue-hotspot-accounting,laravel-queue-service-control,laravel-queue-security,laravel-queue-monitoring,laravel-queue-tenant-management,laravel-queue-user-management,laravel-queue-auth-tracking,laravel-queue-packages
priority=999
