# backend/supervisor/laravel-queue-metrics.conf
[program:laravel-queue-dashboard]
command=/usr/local/bin/php artisan queue:work --queue=dashboard --sleep=1 --tries=3 --timeout=120 --max-time=3600 --memory=128 --backoff=5,15,30
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=15
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:laravel-queue-metrics]
command=/usr/local/bin/php artisan queue:work --queue=metrics --sleep=1 --tries=3 --timeout=120 --max-time=3600 --memory=256 --backoff=5,15,30
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=10
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:laravel-queue-reports]
command=/usr/local/bin/php artisan queue:work --queue=reports --sleep=10 --tries=3 --timeout=600 --max-time=3600 --memory=256 --backoff=30,60,120
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=90
stopsignal=TERM
priority=20
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[group:laravel-queues-metrics]
programs=laravel-queue-dashboard,laravel-queue-metrics,laravel-queue-reports
priority=999
