# backend/supervisor/laravel-queue-router.conf
[program:laravel-queue-router-checks]
command=/usr/local/bin/php artisan queue:work --queue=router-checks --sleep=3 --tries=3 --timeout=120 --max-time=3600 --memory=128 --backoff=5,15,45
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=10
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=false
stdout_logfile=/var/www/html/storage/logs/router-checks-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/router-checks-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-router-data]
command=/usr/local/bin/php artisan queue:work --queue=router-data --sleep=2 --tries=3 --timeout=60 --max-time=1800 --memory=128 --backoff=2,5,15
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=20
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=false
stdout_logfile=/var/www/html/storage/logs/router-data-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/router-data-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-router-provisioning]
command=/usr/local/bin/php artisan queue:work --queue=router-provisioning --sleep=2 --tries=5 --timeout=600 --max-time=1800 --memory=256 --backoff=30,60,120,300,600
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=90
stopsignal=TERM
priority=10
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=false
stdout_logfile=/var/www/html/storage/logs/router-provisioning-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/router-provisioning-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-provisioning]
command=/usr/local/bin/php artisan queue:work --queue=provisioning --sleep=2 --tries=5 --timeout=300 --max-time=1800 --memory=256 --backoff=10,30,60,120,300
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=90
stopsignal=TERM
priority=10
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=false
stdout_logfile=/var/www/html/storage/logs/provisioning-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/provisioning-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-router-monitoring]
command=/usr/local/bin/php artisan queue:work --queue=router-monitoring --sleep=1 --tries=3 --timeout=30 --max-time=3600 --memory=128 --backoff=1,3,5
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=5
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=false
stdout_logfile=/var/www/html/storage/logs/router-monitoring-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/router-monitoring-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-hotspot-provisioning]
command=/usr/local/bin/php artisan queue:work --queue=hotspot-provisioning --sleep=2 --tries=5 --timeout=300 --max-time=1800 --memory=256 --backoff=10,30,60,120,300
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=90
stopsignal=TERM
priority=5
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=false
stdout_logfile=/var/www/html/storage/logs/hotspot-provisioning-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/hotspot-provisioning-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-subscription-reconnection]
command=/usr/local/bin/php artisan queue:work --queue=subscription-reconnection --sleep=2 --tries=3 --timeout=90 --max-time=1800 --memory=128 --backoff=5,15,30
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=10
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=false
stdout_logfile=/var/www/html/storage/logs/subscription-reconnection-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/subscription-reconnection-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-vpn-provisioning]
command=/usr/local/bin/php artisan queue:work --queue=vpn-provisioning --sleep=2 --tries=3 --timeout=120 --max-time=1800 --memory=256 --backoff=10,30,60
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=5
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=false
stdout_logfile=/var/www/html/storage/logs/vpn-provisioning-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/vpn-provisioning-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[group:laravel-queues-router]
programs=laravel-queue-router-checks,laravel-queue-router-data,laravel-queue-router-provisioning,laravel-queue-provisioning,laravel-queue-router-monitoring,laravel-queue-hotspot-provisioning,laravel-queue-subscription-reconnection,laravel-queue-vpn-provisioning
priority=999
