# backend/supervisor/laravel-scheduler.conf
[program:laravel-scheduler]
command=/bin/sh /var/www/html/docker/scheduler-loop.sh
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
numprocs=1
startsecs=5
stopwaitsecs=60
stopsignal=TERM
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3