# backend/supervisor/laravel-scheduler.conf
[program:laravel-scheduler]
command=/usr/local/bin/php /var/www/html/artisan schedule:work
directory=/var/www/html
user=www-data
autostart=true
autorestart=true
numprocs=1
startsecs=0
stopwaitsecs=60
stopsignal=TERM
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3