# backend/supervisor/laravel-queue.conf
[program:laravel-queue-default]
command=/usr/local/bin/php artisan queue:work database --queue=default --sleep=5 --tries=3 --timeout=90 --max-time=3600 --memory=128 --backoff=3,10,30
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=5
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/default-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/default-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-router-checks]
command=/usr/local/bin/php artisan queue:work database --queue=router-checks --sleep=3 --tries=3 --timeout=120 --max-time=3600 --memory=128 --backoff=5,15,45
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=10
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/router-checks-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/router-checks-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-router-data]
command=/usr/local/bin/php artisan queue:work database --queue=router-data --sleep=2 --tries=3 --timeout=60 --max-time=1800 --memory=128 --backoff=2,5,15
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=20
user=www-data
numprocs=4
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/router-data-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/router-data-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-log-rotation]
command=/usr/local/bin/php artisan queue:work database --queue=log-rotation --sleep=30 --tries=3 --timeout=120 --max-time=3600 --memory=128
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=90
stopsignal=TERM
priority=30
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/log-rotation-queue.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/log-rotation-queue-error.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-payments]
command=/usr/local/bin/php artisan queue:work database --queue=payments --sleep=2 --tries=3 --timeout=120 --max-time=1800 --memory=128 --backoff=5,15,60
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=120
stopsignal=TERM
priority=5
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/payments-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/payments-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-payment-checks]
command=/usr/local/bin/php artisan queue:work database --queue=payment-checks --sleep=2 --tries=3 --timeout=90 --max-time=1800 --memory=128 --backoff=5,15,60
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=5
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/payment-checks-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/payment-checks-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-router-provisioning]
command=/usr/local/bin/php artisan queue:work database --queue=router-provisioning --sleep=2 --tries=5 --timeout=600 --max-time=1800 --memory=256 --backoff=30,60,120,300,600
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=90
stopsignal=TERM
priority=10
user=www-data
numprocs=3
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/router-provisioning-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/router-provisioning-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-dashboard]
command=/usr/local/bin/php artisan queue:work database --queue=dashboard --sleep=3 --tries=3 --timeout=120 --max-time=3600 --memory=128 --backoff=5,15,30
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=15
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/dashboard-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/dashboard-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-hotspot-sms]
command=/usr/local/bin/php artisan queue:work database --queue=hotspot-sms --sleep=2 --tries=3 --timeout=45 --max-time=1800 --memory=128 --backoff=3,10,30
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=5
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/hotspot-sms-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/hotspot-sms-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-hotspot-sessions]
command=/usr/local/bin/php artisan queue:work database --queue=hotspot-sessions --sleep=3 --tries=3 --timeout=60 --max-time=1800 --memory=128 --backoff=3,10,30
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=10
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/hotspot-sessions-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/hotspot-sessions-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-hotspot-accounting]
command=/usr/local/bin/php artisan queue:work database --queue=hotspot-accounting --sleep=5 --tries=2 --timeout=120 --max-time=3600 --memory=128 --backoff=10,30,60
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=15
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/hotspot-accounting-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/hotspot-accounting-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-notifications]
command=/usr/local/bin/php artisan queue:work database --queue=notifications --sleep=10 --tries=3 --timeout=60 --max-time=3600 --memory=128 --backoff=5,15,30
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=20
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/notifications-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/notifications-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-service-control]
command=/usr/local/bin/php artisan queue:work database --queue=service-control --sleep=2 --tries=3 --timeout=90 --max-time=1800 --memory=128 --backoff=3,10,30
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=5
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/service-control-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/service-control-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-provisioning]
command=/usr/local/bin/php artisan queue:work database --queue=provisioning --sleep=2 --tries=5 --timeout=300 --max-time=1800 --memory=256 --backoff=10,30,60,120,300
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=90
stopsignal=TERM
priority=10
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/provisioning-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/provisioning-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-router-monitoring]
command=/usr/local/bin/php artisan queue:work database --queue=router-monitoring --sleep=3 --tries=3 --timeout=120 --max-time=3600 --memory=128 --backoff=5,15,45
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=15
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/router-monitoring-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/router-monitoring-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-packages]
command=/usr/local/bin/php artisan queue:work database --queue=packages --sleep=2 --tries=3 --timeout=90 --max-time=3600 --memory=128 --backoff=3,10,30
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=5
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/packages-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/packages-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-broadcasts]
command=/usr/local/bin/php artisan queue:work database --queue=broadcasts --sleep=1 --tries=3 --timeout=30 --max-time=3600 --memory=128 --backoff=2,5,10
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=30
stopsignal=TERM
priority=3
user=www-data
numprocs=3
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/broadcasts-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/broadcasts-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-security]
command=/usr/local/bin/php artisan queue:work database --queue=security --sleep=2 --tries=3 --timeout=60 --max-time=3600 --memory=128 --backoff=3,10,30
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=5
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/security-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/security-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-monitoring]
command=/usr/local/bin/php artisan queue:work database --queue=monitoring --sleep=5 --tries=3 --timeout=300 --max-time=3600 --memory=128 --backoff=10,30,60
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=5
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/monitoring-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/monitoring-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-emails]
command=/usr/local/bin/php artisan queue:work database --queue=emails --sleep=2 --tries=3 --timeout=120 --max-time=3600 --memory=128 --backoff=5,15,30
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=3
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/emails-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/emails-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-tenant-management]
command=/usr/local/bin/php artisan queue:work database --queue=tenant-management --sleep=2 --tries=3 --timeout=180 --max-time=3600 --memory=256 --backoff=5,15,30
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=90
stopsignal=TERM
priority=5
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/tenant-management-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/tenant-management-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-user-management]
command=/usr/local/bin/php artisan queue:work database --queue=user-management --sleep=2 --tries=3 --timeout=120 --max-time=3600 --memory=256 --backoff=5,15,30
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=5
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/user-management-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/user-management-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-auth-tracking]
command=/usr/local/bin/php artisan queue:work database --queue=auth-tracking --sleep=3 --tries=3 --timeout=60 --max-time=3600 --memory=128 --backoff=3,10,30
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=10
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/auth-tracking-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/auth-tracking-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-hotspot-provisioning]
command=/usr/local/bin/php artisan queue:work database --queue=hotspot-provisioning --sleep=2 --tries=5 --timeout=300 --max-time=1800 --memory=256 --backoff=10,30,60,120,300
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=5
startsecs=10
stopwaitsecs=90
stopsignal=TERM
priority=5
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/hotspot-provisioning-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/hotspot-provisioning-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-subscription-reconnection]
command=/usr/local/bin/php artisan queue:work database --queue=subscription-reconnection --sleep=2 --tries=3 --timeout=90 --max-time=1800 --memory=128 --backoff=5,15,30
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=10
user=www-data
numprocs=1
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/subscription-reconnection-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/subscription-reconnection-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[program:laravel-queue-vpn-provisioning]
command=/usr/local/bin/php artisan queue:work database --queue=vpn-provisioning --sleep=2 --tries=3 --timeout=120 --max-time=1800 --memory=256 --backoff=10,30,60
directory=/var/www/html
environment=LARAVEL_ENV="production"
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=5
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/vpn-provisioning-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7
stderr_logfile=/var/www/html/storage/logs/vpn-provisioning-queue-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7
stopasgroup=true
killasgroup=true

[group:laravel-queues]
programs=laravel-queue-default,laravel-queue-router-checks,laravel-queue-router-data,laravel-queue-log-rotation,laravel-queue-payments,laravel-queue-payment-checks,laravel-queue-router-provisioning,laravel-queue-dashboard,laravel-queue-hotspot-sms,laravel-queue-hotspot-sessions,laravel-queue-hotspot-accounting,laravel-queue-notifications,laravel-queue-service-control,laravel-queue-provisioning,laravel-queue-router-monitoring,laravel-queue-packages,laravel-queue-broadcasts,laravel-queue-security,laravel-queue-monitoring,laravel-queue-tenant-management,laravel-queue-user-management,laravel-queue-auth-tracking,laravel-queue-hotspot-provisioning,laravel-queue-subscription-reconnection,laravel-queue-vpn-provisioning
priority=999
